# .config/systemd/user/brightness_lg_1.timer
[Unit]
Description='On *-*-* 07:00:00 run setLGBrightness 10'

[Timer]
Persistent=yes
OnCalendar=*-*-* 07:00:00

[Install]
WantedBy=timers.target

