# .config/systemd/user/brightness_lg_3.timer
[Unit]
Description='On *-*-* 19:25:00 run setLGBrightness 1'

[Timer]
Persistent=yes
OnCalendar=*-*-* 19:25:00

[Install]
WantedBy=timers.target

