# .config/systemd/user/brightness_lg_2.timer
[Unit]
Description='On *-*-* 11:00:00 run brightness_lg_set 30'

[Timer]
Persistent=yes
OnCalendar=*-*-* 11:00:00

[Install]
WantedBy=timers.target

