# .config/systemd/user/catt_default_day.timer
[Unit]
Description='On *-*-* 06:15:00 run cr'

[Timer]
Persistent=yes
OnCalendar=*-*-* 06:15:00

[Install]
WantedBy=timers.target

