# .config/systemd/user/catt_default_day.timer
[Unit]
Description='On *-*-* 03:15:00 run cr'

[Timer]
Persistent=yes
OnCalendar=*-*-* 03:15:00

[Install]
WantedBy=timers.target

