# .config/systemd/user/daily.timer
[Unit]
Description=/usr/bin/fish -c daily

[Timer]
Persistent=yes
OnCalendar=daily
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
