# .config/systemd/user/daily_ffsmallpar.timer
[Unit]
Description='On *-*-* 23:15:00 run ffsmallpar'

[Timer]
Persistent=yes
OnCalendar=*-*-* 23:15:00

[Install]
WantedBy=timers.target

