# .config/systemd/user/dl_monthly.timer
[Unit]
Description=/usr/bin/fish -c dl_monthly

[Timer]
Persistent=yes
OnCalendar=monthly
RandomizedDelaySec=90m

[Install]
WantedBy=timers.target
