# .config/systemd/user/duckdns.service
[Unit]
Description=Duck DNS Updater

[Service]
Type=simple
ExecStart=/home/xk/bin/duckdns

[Install]
WantedBy=default.target
