# .config/systemd/user/firefox.service
[Unit]
Description=firefox

[Service]
Type=simple
RemainAfterExit=no
Environment="DISPLAY=:0"
ExecStart="/home/xk/.local/bin/firefox/firefox"
