# .config/systemd/user/hourly.timer
[Unit]
Description='On hourly run fish -c hourly'

[Timer]
Persistent=yes
OnCalendar=hourly

[Install]
WantedBy=timers.target

