# .config/systemd/user/magnets.timer
[Unit]
Description='On daily run magnets'

[Timer]
Persistent=yes
OnCalendar=daily

[Install]
WantedBy=timers.target

