# .config/systemd/user/midnight.timer
[Unit]
Description='On *-*-* 20:45:00 run midnight'

[Timer]
Persistent=yes
OnCalendar=*-*-* 20:45:00

[Install]
WantedBy=timers.target

