# .config/systemd/user/monthly.service
[Unit]
Description=/usr/bin/fish -c monthly
Requires=mnt-d.mount

[Service]
Type=simple
RemainAfterExit=no
ExecStart=
ExecStart="/usr/bin/fish" "-c" "monthly"
