# .config/systemd/user/noonday.timer
[Unit]
Description='On *-*-* 10:30:00 run noonday'

[Timer]
Persistent=yes
OnCalendar=*-*-* 10:30:00

[Install]
WantedBy=timers.target

