# .config/systemd/user/tabs-cine.service
[Unit]
Description=/usr/bin/fish -c cine
Requires=mnt-d.mount

[Service]
Type=simple
RemainAfterExit=no
Environment="DISPLAY=:0"
ExecStart=
ExecStart="/usr/bin/fish" "-c" "cine"
