# .config/systemd/user/tabs-extra.timer
[Unit]
Description=/usr/bin/fish -c openmctabs /home/xk/mc/tabs

[Timer]
Persistent=yes
OnCalendar=*-*-* 17:05

[Install]
WantedBy=timers.target
