# .config/systemd/user/tabs-music.service
[Unit]
Description=/usr/bin/fish -c opentabs /home/xk/mc/_music

[Service]
Type=simple
RemainAfterExit=no
Environment="DISPLAY=:0"
ExecStart=
ExecStart="/usr/bin/fish" "-c" "opentabs /home/xk/mc/_music"
