# .config/systemd/user/tabs-extra.timer
[Unit]
Description=/usr/bin/fish -c opentabs /home/xk/mc/tabs.txt

[Timer]
Persistent=yes
OnCalendar=*-*-* 17:05

[Install]
WantedBy=timers.target
