# .config/systemd/user/tabs-music.timer
[Unit]
Description=/usr/bin/fish -c opentabs /home/xk/mc/_music

[Timer]
Persistent=yes
OnCalendar=*-*-* 17:15

[Install]
WantedBy=timers.target
