# .config/autostart/autostart.sh.desktop
[Desktop Entry]
Exec=/home/xk/.config/autostart.sh
Icon=dialog-scripts
Name=autostart.sh
Path=
Type=Application
X-KDE-AutostartScript=true
