# .config/systemd/user/xdg-desktop-portal-gtk.service.d/override.conf
[Service]
Environment=XAUTHORITY=/run/user/1000/Xauthority
